# monitor/crontab
0 9 * * * python /app/monitor/monitor.py
