# collector/crontab
# Run the collector every hour at :00
0 * * * * python /app/collector/collect.py
